# Puppet/modules/old_module/logs/file/logstash_server.conf
start on (filesystem and net-device-up)
stop on runlevel [!2345]
respawn
respawn limit 5 30
limit nofile 65550 65550
chdir /opt/logstash/agent
setuid root
setgid root
script
  export LS_HEAP_SIZE="595M"
  export LOGSTASH_HOME="/opt/logstash/agent"
  export HOME=$LOGSTASH_HOME
  export LOGSTASH_OPTS="agent -f /opt/logstash/agent/etc/conf.d/ -l /opt/logstash/agent/log/logstash.log -w 1"
  exec $LOGSTASH_HOME/bin/logstash $LOGSTASH_OPTS
end script
emits logstash-server-running

